# automation/com.courseworkninja.dailyblog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.courseworkninja.dailyblog</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/run-daily-blog.sh</string>
        <string>3</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>WorkingDirectory</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/logs/launchd-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/logs/launchd-stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
