# automation/com.courseworkninja.telegram.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.courseworkninja.telegram</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/telegram-listener.sh</string>
    </array>
    
    <key>WorkingDirectory</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja</string>
    
    <key>StandardOutPath</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/logs/telegram-stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/mx/Documents/WGU/Claude/coursework-ninja/automation/logs/telegram-stderr.log</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
